# deploy/pinkman-backend.service
# Pinkman backend — isolated systemd service.
# Binds to 127.0.0.1 ONLY, so it is not exposed on the network and requires
# zero changes to Oracle VCN / subnets / security lists / ingress rules.
#
# Install:
#   sudo cp deploy/pinkman-backend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pinkman-backend
#
# Remove completely (rollback):
#   sudo systemctl disable --now pinkman-backend
#   sudo rm /etc/systemd/system/pinkman-backend.service
#   sudo systemctl daemon-reload

[Unit]
Description=Pinkman FastAPI backend (localhost only)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/pinkman/backend
EnvironmentFile=/home/ubuntu/pinkman/backend/.env
# --host 127.0.0.1 keeps it off the network. Change --port only if 8000 is taken.
ExecStart=/home/ubuntu/pinkman/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3
# Light isolation so this service can't disturb the rest of the box.
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
